# stvor-ws-relay.service
[Unit] 
Description=STVOR WebSocket Relay Server
After=network.target

[Service]
Type=simple
User=nodejs
WorkingDirectory=/opt/stvor-api
ExecStart=/usr/bin/node packages/sdk/relay-server.js
Environment=NODE_ENV=production
Environment=PORT=8080
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=stvor-ws-relay

[Install]
WantedBy=multi-user.target